# layers/meta-opentrons/recipes-robot/opentrons-live-stream/files/opentrons-live-stream.service
[Unit]
Requires=nginx.service
After=nginx.service

[Service]
ExecStart=/usr/bin/opentrons-live-stream.sh
Restart=always

[Install]
WantedBy=multi-user.target