# layers/meta-opentrons/recipes-robot/opentrons-systemd-units/files/opentrons-clear-fontconfig-cache.service
[Service]
Type=oneshot
ExecStart=/usr/bin/ot-clear-fontconfig-cache

[Install]
WantedBy=basic.target
