# layers/meta-opentrons/recipes-multimedia/gstreamer-1.0/files/gstd.service
[Unit]
Description=gstd
After=weston.service
Requires=weston.service
Before=multi-user.target

[Service]
Type=forking
ExecStart=/usr/bin/gstd
Environment="DISPLAY=:0:0" "XDG_SESSION_TYPE=wayland" "XDG_SESSION_DESKTOP=kiosk" "XDG_RUNTIME_DIR=/%t/user/%U"

[Install]
WantedBy=basic.target
