# layers/meta-opentrons/recipes-robot/opentrons-pyro-nameserver/files/opentrons-pyro-nameserver.service
[Unit]
Requires=nginx.service
After=nginx.service

[Service]
# Begin the Opentrons Pyro5 Nameserver
ExecStart=/bin/sh -c 'pyro5-ns'   

Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target