# layers/meta-opentrons/recipes-robot/opentrons-robot-app/files/opentrons-robot-app-devtools.socket
[Socket]
ListenStream=9223

[Install]
WantedBy=sockets.target
