# layers/meta-opentrons/recipes-robot/opentrons-robot-app/files/opentrons-robot-app.service.in
[Unit]
Description=Start the Opentrons Robot App
After=weston@root.service
Requires=weston@root.service

[Service]
Type=notify
Environment=@@application-environment@@
ExecStart=/usr/bin/opentrons-robot-app.sh
ExecStartPost=/bin/systemctl stop opentrons-loading.service
ExecStopPost=/bin/systemctl start opentrons-loading.service
TimeoutStartSec=30
ExitType=cgroup
NotifyAccess=all

[Install]
WantedBy=multi-user.target
